# henzai-daemon/systemd/henzai-daemon.service
[Unit]
Description=henzai Daemon - AI-First GNOME Desktop Assistant
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/henzai-daemon
Restart=on-failure
RestartSec=5

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target










